# docker/services.d/sshd/run
#!/bin/sh
mkdir -p /run/sshd
/usr/sbin/sshd -t &&
exec /usr/sbin/sshd -D -e
